# crontab
*/5 * * * * bundle exec rake /rails/lib/tasks/create_hellos.rake create_hellos
